# nis.service
[Unit]
Description=NIS Service

[Service]
ExecStart=root/package/nix.runNis.sh
Restart=always

[Install]
WantedBy=multi-user.target
