# supernode.service
[Unit]
Description=Supernode Service

[Service]
ExecStart=/root/supernodeStart.sh
Restart=always

[Install]
WantedBy=multi-user.target
